# deploy/malkovich.service
[Unit]
Description=Malkovich platform service
After=network.target

[Service]
Type=simple
User=garage44
Group=garage44
WorkingDirectory=/home/garage44/garage44/packages/malkovich
Environment="NODE_ENV=production"
Environment="BUN_ENV=production"
Environment="PATH=/home/garage44/.bun/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="WEBHOOK_SECRET=your-webhook-secret-here"
ExecStart=/home/garage44/.bun/bin/bun run server -- --port 3032
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
