# deploy/pr-cleanup.service
[Unit]
Description=Cleanup stale PR deployments
After=network.target

[Service]
Type=oneshot
User=garage44
Group=garage44
WorkingDirectory=/home/garage44/garage44/packages/malkovich
Environment="PATH=/home/garage44/.bun/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/home/garage44/.bun/bin/bun run malkovich cleanup-stale-prs
StandardOutput=journal
StandardError=journal
